# systemd/user/niri-eww.service
[Unit]
Description=Eww bar for Niri
PartOf=niri.service
After=niri.service

[Service]
ExecStartPre=-/usr/bin/eww close bar 2>/dev/null
ExecStart=/usr/bin/eww --no-daemonize daemon
ExecStartPost=/bin/sh -c 'for i in $(seq 1 50); do eww ping >/dev/null 2>&1 && break; sleep 0.1; done; eww open --id primary bar'
ExecStopPost=-/usr/bin/eww close-all
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=niri.service
